# rootfs/etc/init.d/tty
#!/sbin/openrc-run
start ()
{
  ebegin "Init tty.."

  # if ! [ /dev/tty1 -ef /dev/console ]; then
      # lxc_console="console"
  # fi

  ttys="tty1 ttyS0 ttyACM0 ttyAMA0 hvc0"
  inittab_file="/etc/inittab"

  for tty in $ttys; do
    if [ -e /dev/$tty ]; then
      if ! grep -q "^$tty:" "$inittab_file"; then
        case "$tty" in
          console)
            echo "console::respawn:/sbin/getty 38400 $tty" >> "$inittab_file"
            ;;
          *)
            echo "$tty::respawn:/sbin/getty 38400 $tty" >> "$inittab_file"
            ;;
        esac
      fi
    fi
  done

  kill -HUP 1

  eend 0
}

