# rootfs/etc/local.d/podman-socket.start
#!/bin/sh

SOCKET=/run/podman/podman.sock
if [ -e "$SOCKET" ]; then
  chown root:wheel "$SOCKET"
  chmod 660       "$SOCKET"
fi
