# rootfs/etc/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
::shutdown:/sbin/openrc shutdown

console::respawn:/sbin/getty -L console 115200 vt100

::ctrlaltdel:/sbin/reboot

